# package/debian/DEBIAN/prerm
#!/bin/sh

/etc/init.d/raspwave stop
/usr/sbin/update-rc.d -f raspwave remove
