# package/postinst
#!/bin/sh

/usr/sbin/update-rc.d raspwave defaults
/bin/sleep 5
/etc/init.d/raspwave start
